# systemd/flatpak-update.service
[Unit]
Description=Update Flatpak applications

[Service]
Type=oneshot
ExecStart=/usr/bin/flatpak update -y
